# conductor.service
[Unit]
Description=Conductor Web Application Server

[Service]
User=conductor
Type=forking
ExecStart=
ExecStop=


